# App/scripts/spec-extraction-worker.service
[Unit]
Description=Spec Extraction Worker Service
After=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/spec-extraction.env
WorkingDirectory=/opt/spec-extraction
ExecStart=/bin/bash /opt/spec-extraction/App/scripts/run_worker.sh
Restart=always
RestartSec=5
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target
